# launch/com.imdiscobridge.responder.plist
<?xml version="1.0" encoding="UTF-8"?>
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.imdiscobridge.responder</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOURUSERNAME/imdiscobridge/venv/bin/python</string>
        <string>/Users/YOURNAME/imdiscobridge/scripts/responder.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/imdiscobridge_responder.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/imdiscobridge_responder.err</string>
</dict>
</plist>
